# deploy/llm-gateway.service
[Unit]
Description=llm-gateway FastAPI service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=llmgateway
Group=llmgateway
WorkingDirectory=/opt/llm-gateway
EnvironmentFile=/opt/llm-gateway/.env
ExecStart=/opt/llm-gateway/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
ExecStartPost=/opt/llm-gateway/scripts/check_ready.sh
Restart=always
RestartSec=3
TimeoutStartSec=90
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
